# crontab
0 10 * * * /usr/bin/bash /opt/monero/.local/src/srv-xmr-nodes/utils.sh >/dev/null 2>&1
