# etc/systemd/system/monerod-testnet.service
[Unit]
Description=Monero testnet Daemon
After=network.target

[Service]
User=monero
Type=forking
GuessMainPID=no
ExecStart=/opt/monero/.local/bin/monerod \
    --detach \
    --testnet \
    --config-file=/etc/monero/testnet.conf
Restart=always

[Install]
WantedBy=multi-user.target

# vim: filetype=systemd ts=4 sw=4 et
